# mkosi.images/system/mkosi.extra/usr/lib/systemd/system/chronyd.service.d/dirs.conf
[Service]
#ReadWritePaths=
StateDirectory=chrony
StateDirectoryMode=0750
RuntimeDirectory=chrony
RuntimeDirectoryMode=0750
LogsDirectory=chrony
LogsDirectoryMode=0750
#User=chrony
#Group=chrony
#DynamicUser=yes
#Environment="OPTIONS=-U"
ExecStartPre=-chown chrony:chrony "%t/chrony"
