# mkosi.images/system/mkosi.extra/usr/lib/systemd/system/chronyd.service.d/dirs.conf
[Service]
StateDirectory=chrony
StateDirectoryMode=0750
RuntimeDirectory=chrony
RuntimeDirectoryMode=0750
