# mkosi.profiles/debug/mkosi.conf.d/90-sshd-vsock/mkosi.extra/usr/lib/systemd/system/ssh@.service.d/use-pam.conf
[Service]
# Remove the original ExecStart=sshd -i -o UsePAM=no by assigning an empty ExecStart=
ExecStart=
ExecStart=sshd -i
