# include/systemd-udev/mkosi.extra/usr/lib/systemd/system/systemd-bless-boot.service.d/mount-boot.conf
[Unit]
After=boot.automount boot.mount
After=efi.automount efi.mount
Wants=boot.automount
