# include/issuegen/mkosi.extra/usr/lib/systemd/system/issuegen-netif@.service
[Unit]
Description=Write IP address template for interface %I to %t/issue.d/
Documentation=man:issue(5) man:agetty(8) man:ip(8) man:ip-link(8)
After=network.target

[Service]
Type=oneshot
RuntimeDirectory=issue.d
RuntimeDirectoryPreserve=yes
RuntimeDirectoryMode=0755
ExecStart=/usr/libexec/issuegen-netif '%t/issue.d/' '%I' --nomaster
ExecStartPost=-agetty --reload

[Install]
WantedBy=multi-user.target
