# include/state-dirs/mkosi.extra/usr/lib/systemd/system/sssd-kcm.service.d/dirs.conf
[Service]
LogsDirectory=sssd
LogsDirectoryMode=0750
