# include/state-dirs/mkosi.extra/usr/lib/systemd/system/ldconfig.service.d/dirs.conf
[Service]
CacheDirectory=ldconfig
CacheDirectoryMode=0700
