# mkosi.images/system/mkosi.conf.d/60-readonly-fs/mkosi.extra/usr/lib/systemd/system/tmpfiles-create@.service
[Unit]
DefaultDependencies=no
RequiresMountsFor=%f
Before=shutdown.target
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart=systemd-tmpfiles --create --prefix='%f'
