# include/state-dirs/mkosi.extra/usr/lib/systemd/system/sssd.service.d/dirs.conf
[Service]
LogsDirectory=%p
LogsDirectoryMode=0750
