# mkosi.images/base/mkosi.conf.d/44-libvirt/mkosi.extra/usr/lib/systemd/system/virtqemud.service.d/90-allow-jit-memory.conf
[Service]
# qemu does not work without executable writable memory
MemoryDenyWriteExecute=no
