# mkosi.images/system/mkosi.extra/usr/lib/systemd/system/setroubleshoot.service.d/dirs.conf
[Service]
StateDirectory=setroubleshoot
StateDirectoryMode=0700
RuntimeDirectory=setroubleshoot
RuntimeDirectoryMode=0711
