# include/state-dirs/mkosi.extra/usr/lib/systemd/system/gssproxy.service.d/dirs.conf
[Service]
StateDirectory=gssproxy/clients gssproxy/rcache
StateDirectoryMode=0700
