# include/state-dirs/mkosi.extra/usr/lib/systemd/system/udisks2.service.d/dirs.conf
[Service]
StateDirectory=%p
StateDirectoryMode=0750
RuntimeDirectory=%p
RuntimeDirectoryMode=0750
