# mkosi.images/system/mkosi.extra/usr/lib/systemd/system/systemd-journald.service.d/logdir.conf
[Service]
LogsDirectory=journal
LogsDirectoryMode=2755
